# contrib/womo/etc/dnsmasq/systemd/lan-mode.timer
[Unit]
Description=Check WAN status periodically

[Timer]
OnBootSec=10
OnUnitActiveSec=10

[Install]
WantedBy=timers.target
